# run/supervisord/celery.conf
[program:celery]
command=
autostart=true
autorestart=true
stderr_logfile=/var/log/flower.err.log
stdout_logfile=/var/log/flower.out.log
